# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/app/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 7070
directory=/app/backend
autostart=true
autorestart=true
environment=DATABASE_URL="sqlite:////data/laan.db",HOST_AGENT_URL="http://host.docker.internal:7072"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
